# wlan0-check.timer
[Unit]
Description=Run wlan0-check.service ever 2 minutes

[Timer]
OnCalendar=*:0/02
